# package/DEBIAN/postinst
echo "......................................................"
echo "..............BroadCast_deb installing...................."
echo "......................................................"

systemctl enable BroadCast
systemctl start BroadCast