# package/DEBIAN/postrm
if [ "$1" = "upgrade" ] ; then
   echo "upgrade"
elif [ "$1" = "remove" ] || [ "$1" = "purge" ] ; then
   echo "remove"
   systemctl disable BroadCast
fi